# dummy-systemd-service/dummy.service
[Unit]
Description=Dummy script
After=network.target

[Service]
ExecStart=/bin/bash /opt/dummy/dummy.sh
Restart=always